# systemd/system/kubuntu-lldp-agent.service
[Unit]
Description=KubuntuLLDP privileged network agent
After=network-pre.target
Before=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kubuntu-lldp-agent
Restart=on-failure
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target

